# doom-fire-wallpaper.service
[Unit]
Description=DOOM Fire Wallpaper for Hyprpaper
After=graphical-session.target hyprpaper.service

[Service]
Type=simple
ExecStart=/usr/bin/doom-fire-wallpaper
Restart=on-failure
Nice=10

[Install]
WantedBy=default.target